# LFS_11.2_systemd/etc/profile.d/F99-colorterm.sh
# set the COLORTERM variable if in an appropriate terminal
case "$TERM" in 
    xterm-kitty) export COLORTERM=fullcolor ;;
esac
