# LFS_11.2_systemd/etc/profile.d/Z99-motd.sh
if [ "$SHLVL" -le 1 ]; then . /etc/motd.sh; fi
