# etc/profile.d/meson.sh
export MESON_ARGS='-Dlibdir=/usr/lib -Dbackend=ninja -Dstrip=true -Ddebug=false'
