# etc/profile.d/rustc.sh
pathprepend /opt/rustc/bin           PATH
